# bin/systemd/check-battery-user.service
[Unit]
Description=check-battery

[Service]
ExecStart=/home/gaz/bin/battery-watch
