# bin/systemd/suspend@.service
[Unit]
Description=User suspend actions
Before=suspend.target 

[Service]
User=%I
Type=forking
Environment=DISPLAY=:0
ExecStart=/home/gaz/bin/i3lock-color
ExecStartPost=/usr/bin/sleep 1

[Install]
WantedBy=suspend.target
